# deploy/irondag-node.service
[Unit]
Description=IronDAG Blockchain Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/irondag/irondag-blockchain
ExecStartPre=-/usr/bin/fuser -k 8080/tcp
ExecStart=/opt/irondag/irondag-blockchain/target/release/node --port 8080 --rpc-port 8546 --mine --cors-origin "*" --rpc-no-auth
Restart=always
RestartSec=5
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
